# rc.local
#!/bin/bash
# FUNCTION: execute scripts at startup (example file for myscript.sh)
# ARGUMENTS: none
# AUTHOR: jlanssie@gmail.com
# COPYRIGHT: 2021 GNU
# REQUIRES: save in path "/etc/rc.local"

/usr/local/bin/myscript.sh

/etc/sysctl.d
/etc/init.d/procps restart

exit 0