# rc.local
#!/bin/bash
#
# Template for rc.local

/usr/local/bin/myscript.sh

/etc/sysctl.d
/etc/init.d/procps restart

exit 0